# systemd/user/dropbox.service
[Unit]
Description=Dropbox Daemon
After=alopex.service

[Service]
ExecStart=/usr/bin/dropboxd

[Install]
WantedBy=initrc.target
