# systemd/user/dropbox.service
[Unit]
Description=Dropbox Daemon

[Service]
ExecStart=/usr/bin/dropboxd

[Install]
WantedBy=initrc.target
