# systemd/user/themes.service
[Unit]
Description=Normalize Graphical Style
After=xorg.target

[Service]
ExecStart=/usr/bin/xrdb %h/.Xresources

[Install]
WantedBy=initrc.target
